# unified-platform/infra/ansible/roles/self_healing/templates/cluster-watchdog.service.j2
[Unit]
Description=Cluster Self-Healing Watchdog
After=network.target

[Service]
Type=oneshot
ExecStart={{ cluster_root }}/self_healing/systemd_watchdog.sh
StandardOutput=journal
StandardError=journal
